# ops/systemd/vpn-manager.service
[Unit]
Description=DownloaderPro VPN Manager
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/downloaderpro/vpn.env
EnvironmentFile=-/opt/downloaderpro/.env
ExecStart=/opt/downloaderpro/ops/scripts/connect-country.sh
ExecStop=/opt/downloaderpro/ops/scripts/disconnect-vpn.sh
WorkingDirectory=/opt/downloaderpro
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
